# backups/sigma-pipeline.service
[Unit]
Description=SIGMA Master Pipeline - 5 activos x 4 TFs continuo
After=network.target sigma-trainer.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sigma
ExecStart=/opt/sigma_env/bin/python -u engine/live/master_pipeline.py
Restart=always
RestartSec=30
KillMode=control-group
TimeoutStopSec=15
StandardOutput=append:/opt/sigma/results/reports/master_pipeline.log
StandardError=append:/opt/sigma/results/reports/master_pipeline.log

[Install]
WantedBy=multi-user.target
