# backups/sigma-paper-trader.service
[Unit]
Description=SIGMA Paper Trader Daemon
After=sigma-web.service
Wants=sigma-web.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sigma
ExecStart=/opt/sigma_env/bin/python -u /opt/sigma/engine/live/paper_trader_daemon.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/sigma/results/reports/paper_trader_daemon.log
StandardError=append:/opt/sigma/results/reports/paper_trader_daemon.log

[Install]
WantedBy=multi-user.target
