# backups/sigma-frontend.service
[Unit]
Description=SIGMA Web - Next.js Frontend (squantdesk.com)
After=network.target sigma-web.service

[Service]
Type=simple
WorkingDirectory=/opt/sigma_web
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/opt/sigma_web/.env.production
ExecStart=/usr/bin/npm run start
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/sigma_web/server.log
StandardError=append:/opt/sigma_web/server.err.log

[Install]
WantedBy=multi-user.target
