# backups/sigma-web.service
[Unit]
Description=SIGMA Web Dashboard
After=network.target
StartLimitIntervalSec=120
StartLimitBurst=10

[Service]
Environment=TZ=America/Santiago
Type=simple
Restart=always
RestartSec=5
TimeoutStopSec=30
KillMode=mixed
User=root
WorkingDirectory=/opt/sigma
ExecStart=/opt/sigma_env/bin/python /opt/sigma/web_server.py
StandardOutput=append:/opt/sigma/results/reports/web.log
StandardError=append:/opt/sigma/results/reports/web.log

[Install]
WantedBy=multi-user.target
