# backups/sigma-champion-watcher.service
[Unit]
Description=SIGMA Champion Watcher Daemon - reactivo a cambios en models/
After=network.target sigma-web.service
Wants=sigma-web.service

[Service]
Type=simple
WorkingDirectory=/opt/sigma
ExecStart=/opt/sigma_env/bin/python /opt/sigma/champion_watcher_daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/sigma/results/reports/champion_watcher_daemon.log
StandardError=append:/opt/sigma/results/reports/champion_watcher_daemon.log

[Install]
WantedBy=multi-user.target
