# backups/sigma-nextjs.service
[Unit]
Description=SIGMA Web Next.js App
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/sigma_web
EnvironmentFile=/opt/sigma_web/.env.production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=append:/var/log/sigma-nextjs.log
StandardError=append:/var/log/sigma-nextjs.log

[Install]
WantedBy=multi-user.target

